# ccs-backup/systemd/ccs-backup.service.template
[Unit]
Description=CCS durable usage sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__WORKDIR__
ExecStart=__BUN__ run src/cli.ts sync
Nice=10

[Install]
WantedBy=default.target
